# seedbank.service
[Unit]
Description=Seedbank – lokaler semantischer Speicher (OpenClaw)
After=network.target

[Service]
Type=simple
# Pfad anpassen, falls das Repo woanders liegt:
WorkingDirectory=/home/jarvis/seedbank
EnvironmentFile=/home/jarvis/seedbank/.env
ExecStart=/home/jarvis/seedbank/seedbank
StandardOutput=journal
StandardError=journal
SyslogIdentifier=seedbank
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
